# device_app.service
[Unit]
Description=Device Management Flask App
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/my_device_app
ExecStart=/home/pi/my_device_app/venv/bin/python /home/pi/my_device_app/app.py
Restart=always
RestartSec=10
Environment="PATH=/home/pi/my_device_app/venv/bin"

[Install]
WantedBy=multi-user.target